# debian/landscape-client.postrm
#!/bin/sh
set -e

test "$1" = "purge" || exit 0

update-rc.d landscape-client remove >/dev/null || exit $?
getent passwd landscape >/dev/null || exit 0
! deluser --version >/dev/null 2>&1 || exec deluser landscape
echo 'deluser program not available, not removing system user "landscape".' >&2
