# debian/landscape-client.cron.hourly
#!/bin/sh

[ -f "/etc/default/landscape-client" ] || exit 0

. /etc/default/landscape-client

if [ "$RUN" -eq "1" ] 
then
	[ -x /usr/share/smart/smart ] && /usr/share/smart/smart update > /dev/null 2>&1
fi

exit 0
